# autostart/chatterino.desktop
[Desktop Entry]
Exec=chatterino
Name=chatterino
OnlyShowIn=LXQt;
Type=Application
Version=1.0
