# autostart/cmst.desktop
[Desktop Entry]
Exec=cmst
Name=cmst
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
