# autostart/CopyQ.desktop
[Desktop Entry]
Exec=copyq
Name=CopyQ
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
