# autostart/Create OBS Audio Sinks.desktop
[Desktop Entry]
Exec=sh /iao/Executables/Pulseaudio/create-obs-audio-sinks.sh
Hidden=true
Name=Create OBS Audio Sinks
OnlyShowIn=LXQt;
Type=Application
Version=1.0
