# autostart/discord.desktop
[Desktop Entry]
Exec=discordcanary
Name=discordcanary
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
