# autostart/featherpad.desktop
[Desktop Entry]
Exec=sh ~/Executables/Autostart/featherpad-autostart.sh
Name=featherpad
OnlyShowIn=LXQt;
Type=Application
Version=1.0
