# autostart/Firefox.desktop
[Desktop Entry]
Exec=firefox
Name=Firefox
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
