# autostart/lxqt-config.desktop
[Desktop Entry]
Exec=lxqt-config
Name=lxqt-config
OnlyShowIn=LXQt;
Type=Application
Version=1.0
