# autostart/pavucontrol-qt.desktop
[Desktop Entry]
Exec=pavucontrol-qt
Name=pavucontrol-qt
OnlyShowIn=LXQt;
Type=Application
Version=1.0
