# autostart/picom.desktop
[Desktop Entry]
Exec=picom -b --experimental-backends
Name=picom
OnlyShowIn=LXQt;
Type=Application
Version=1.0
