# autostart/qterminal.desktop
[Desktop Entry]
Exec=qterminal
Name=qterminal
OnlyShowIn=LXQt;
Type=Application
Version=1.0
