# autostart/Ryzen Stabilizator.desktop
[Desktop Entry]
Exec=sudo ryzen-stabilizator
Name=Ryzen Stabilizator
Type=Application
Version=1.0
