# autostart/sxhkd.desktop
[Desktop Entry]
Exec=sxhkd
Name=sxhkd
OnlyShowIn=LXQt;
Type=Application
Version=1.0
