# autostart/tint2 Left.desktop
[Desktop Entry]
Exec=tint2 -c ~/Configuration/tint2/iao_left.tint2rc
Name=tint2 Left
OnlyShowIn=LXQt;
Type=Application
Version=1.0
