# autostart/tint2 Right.desktop
[Desktop Entry]
Exec=tint2 -c ~/Configuration/tint2/iao_right.tint2rc
Name=tint2 Right
OnlyShowIn=LXQt;
Type=Application
Version=1.0
